# systemd/ai-assistant-watchdog.service
[Unit]
Description=AI Assistant Watchdog
After=network.target
# Start after bot, but don't require it (watchdog should start even if bot is dead)
Wants=ai-assistant-bot.service

[Service]
Type=simple
User=ddarji
Group=ddarji
WorkingDirectory=/media/ddarji/storage/ai-assistant
ExecStart=/home/ddarji/.nvm/versions/node/v20.20.0/bin/node src/watchdog/watchdog.js
# systemd watches the watchdog — always restart it
Restart=always
RestartSec=10
StandardOutput=append:/media/ddarji/storage/ai-assistant/watchdog.log
StandardError=append:/media/ddarji/storage/ai-assistant/watchdog.log
Environment=NODE_ENV=production
Environment=PATH=/home/ddarji/.nvm/versions/node/v20.20.0/bin:/home/ddarji/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
